# config/autostart/cbatticon.desktop
[Desktop Entry]
Exec=/usr/bin/cbatticon
Name=cbatticon
Comment=a lightweight and fast battery icon that sits in your system tray
Type=Application
