# config/autostart/sxhkd.desktop
[Desktop Entry]
Exec=/usr/bin/sxhkd
Name=sxhkd
Comment=Simple X hotkey daemon
Type=Application
